# debian-rinetd/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:rinetd]
command=/usr/sbin/rinetd -f -c /etc/rinetd.conf